# tags/OSF-1.0/osf/osf-web/debian/cron.d
# should run after logrotate
45 6 * * *	root	/usr/sbin/postfix-logreport-update /var/lib/osf/postfix-logreport/
